# .config/systemd/user/background-changer.timer
[Unit]
Description=Background changer with feh

[Timer]
OnBootSec=1m
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target